# workspace/crontab/laradock
* * * * * laradock php /var/www/yello/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock php /var/www/marbec/artisan schedule:run >> /dev/null 2>&1
